# services/qbittorrent/finish
#!/data/data/com.termux/files/usr/bin/sh
# Kill any orphan qbittorrent-nox process before the next supervised start.
# SIGKILL (-9) is used to guarantee immediate termination; SIGTERM may
# leave the process alive long enough for the next run to fail on bind.
pkill -9 -x qbittorrent-nox 2>/dev/null || true
