# services/audiobookshelf/finish
#!/data/data/com.termux/files/usr/bin/sh
pkill -9 -f "audiobookshelf/index.js" 2>/dev/null || true
pkill -9 -f "proot-distro" 2>/dev/null || true
