# services/adguard/finish
#!/data/data/com.termux/files/usr/bin/sh
# Runit finish script for adguard.
# Called by runsv every time the run script exits (crash, sv stop, sv restart).
#
# AdGuardHome runs as root via "su -c", so SIGTERM sent by runsv to the su
# wrapper process may not propagate to the root child. This script ensures
# any lingering root AdGuardHome process is killed before the next start,
# preventing the sessions.db lock timeout on restart.
su -c 'pkill -x AdGuardHome 2>/dev/null; true'
