# services/beszel-agent/finish
#!/data/data/com.termux/files/usr/bin/sh
# beszel-agent runs as root via su -c; kill root child before next start.
su -c 'pkill -x beszel-agent 2>/dev/null; true'
