# services/beszel-agent/run
#!/data/data/com.termux/files/usr/bin/sh
set -a
. "$HOME/phone-server/.secrets/boot.env"
. "$HOME/phone-server/.secrets/beszel.env"
set +a

sv check $HOME/phone-server/services/beszel || exit 1

# Data for external ssd in not collected (might be Android limitations)
# export EXTRA_FILESYSTEMS="/dev/block/sdg1"

exec su -c "KEY=\"$BESZEL_AGENT_KEY\" HUB_URL=$BESZEL_URL FILESYSTEM=/data $BIN_DIR/beszel-agent >> $LOG_DIR/beszel-agent.log 2>&1"
