# services/webdav/finish
#!/data/data/com.termux/files/usr/bin/sh
# Kill any orphan rclone webdav process.
pkill -x rclone 2>/dev/null || true
