# services/gatus/finish
#!/data/data/com.termux/files/usr/bin/sh
# Kill any orphan gatus process holding port 8080.
pkill -x gatus 2>/dev/null || true
