# services/gatus/run
#!/data/data/com.termux/files/usr/bin/sh
. "$HOME/phone-server/.secrets/boot.env"
. "$HOME/phone-server/.secrets/gatus.env"
sv check AdGuardHome || exit 1

exec env GATUS_CONFIG_PATH="$GATUS_CONFIG" "$BIN_DIR/gatus" >> $LOG_DIR/gatus.log 2>&1
