# services/caddy/finish
#!/data/data/com.termux/files/usr/bin/sh
# caddy runs as root via su -c; kill root child with SIGKILL.
su -c 'pkill -9 -x caddy 2>/dev/null; true'
