# services/caddy/finish
#!/data/data/com.termux/files/usr/bin/sh
# caddy runs as root via su -c; kill root child before next start.
su -c 'pkill -x caddy 2>/dev/null; true'
