# services/beszel/finish
#!/data/data/com.termux/files/usr/bin/sh
# Kill any orphan beszel hub process holding port 8090.
pkill -x beszel 2>/dev/null || true
