# services/beszel/run
#!/data/data/com.termux/files/usr/bin/sh
. $HOME/boot.env
exec $BESZEL_DIR/beszel serve --http "0.0.0.0:$BESZEL_PORT" --dir $BESZEL_DATA  >> $LOG_DIR/beszel-hub.log 2>&1
