# services/beszel/run
#!/data/data/com.termux/files/usr/bin/sh
. "$HOME/phone-server/.secrets/boot.env"
exec $BIN_DIR/beszel serve --http "0.0.0.0:$BESZEL_PORT" --dir $BESZEL_DATA  >> $LOG_DIR/beszel-hub.log 2>&1
